# tmp/droid-config.files
%defattr(-,root,root,-)
/var/lib/environment/ofono/noplugin.conf
/var/lib/environment/nemo/99-hybris.conf
/var/lib/environment/nemo/60-multimedia.conf
/lib/systemd/system/adbd-prepare.service
/lib/systemd/system/graphical.target.wants/droid-late-start.service
/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/lib/systemd/system/droid-battery-monitor.service
/lib/systemd/system/droid-late-start.service
/lib/systemd/system/basic.target.wants/droid-battery-monitor.service
/lib/systemd/system/adbd.service
/lib/systemd/system/droid-hal-init.service
/lib/udev/rules.d/998-droid-system.rules
/lib/udev/rules.d/997-droid-drm.rules
/lib/udev/platform-device
/etc/xdg/QtProject/Sensors.conf
/etc/sensorfw/primaryuse.conf
/etc/sysctl.d/ipv4-settings.conf
/etc/sysctl.d/ipv6-settings.conf
/etc/modprobe.d/blacklist-evbug.conf
/usr/share/environment/nemo/hybris.conf
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/get_partition.sh
/usr/bin/droid/droid-init-done.sh
/usr/bin/droid/kill-cgroup.sh
/usr/sbin/flash-partition
/usr/sbin/adbd-functionfs.sh
/usr/lib/tmpfiles.d/droid-hal.conf
/var/lib/environment/nemo/60-hlte-vibra.conf
/var/lib/environment/compositor/droid-hal-device.conf
/etc/xdg/gstomx.conf
/etc/hw-release
/usr/share/doc/SailfishOS/extended-packagelist-droid
/usr/share/doc/SailfishOS/droid-hal-version-hlte
/usr/share/doc/SailfishOS/hw-release
/usr/share/ngfd/plugins.d/ngf-vibra-hlte.ini
/usr/share/sailfish-version/packagelist.d/droid-hal-version-hlte
/usr/lib/qt5/plugins/feedback/qtfeedback-hlte.ini
/usr/include/private.h
/usr/include/droidmediaconvert.h
/usr/include/droidmediacodec.h
/usr/include/allocator.h
/usr/include/droidmediacamera.h
/usr/include/droidmediabuffer.h
/usr/include/droidmedia.h
/usr/include/droidmediaconstants.h
